# mopidy/etc/services.d/crond/run
#!/usr/bin/execlineb -P

crond -b -L /var/log/cron.log && tail -f /var/log/cron.log
