# music-server/root/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bash

. /etc/default/snapserver

exec \
        /usr/local/bin/snapserver $SNAPSERVER_OPTS
